# deploy/robot-buddy-supervisor.service
[Unit]
Description=Robot Buddy Supervisor
Documentation=https://github.com/0mdb/robot-buddy
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=pi
Group=pi
WorkingDirectory=/home/pi/robot-buddy

# Load runtime flags from /etc/robot-buddy/supervisor.env
# The leading dash makes the file optional (no failure if missing)
EnvironmentFile=-/etc/robot-buddy/supervisor.env

ExecStart=/home/pi/robot-buddy/supervisor/.venv/bin/python -m supervisor $SUPERVISOR_ARGS

# Restart on non-zero exit; give hardware time to settle after a crash
Restart=on-failure
RestartSec=5s

# Cap restart attempts: if it fails 5 times in 60 s, stop trying
StartLimitIntervalSec=60
StartLimitBurst=5

# Log to journald (view with: journalctl -fu robot-buddy-supervisor)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=robot-buddy

# Serial ports (dialout) and camera (video) access
SupplementaryGroups=dialout video

# Give the process a clean shutdown window before SIGKILL
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
